# openwrt-addons/etc/init.d/netfilter
#!/bin/sh /etc/rc.common

START=45

start() {
    . /tmp/loader
	_netfilter start
}

stop() {
    . /tmp/loader
	_netfilter stop
}
