# openwrt-addons/etc/init.d/S41build_static_netparam
#!/bin/sh /etc/rc.common

START=41

EXTRA_COMMANDS="call"
EXTRA_HELP=<<EOF
        call same as start
EOF



call()
{
	start
}

start()
{
	. /tmp/loader
	_netparam check
}
