# openwrt-addons/etc/init.d/S41build_static_netparam
#!/bin/sh /etc/rc.common

START=41

call()
{
	start
}

start()
{
	. /tmp/loader
	_netparam check
}
