# openwrt-addons/etc/crontabs/root
* * * * * . /usr/sbin/cron.minutely; . /tmp/loader; _watch switch_linkstatus
