# files/etc/uci-defaults/11-migrate-meshwizard.sh
#!/bin/sh 

. /lib/functions.sh 


config_load meshwizard || return 0
# vim: set filetype=sh ai noet ts=4 sw=4 sts=4 :
