# files/etc/uci-defaults/11-migrate-meshwizard.sh
#!/bin/sh 

. /lib/functions.sh 


config_load meshwizard || return 0
