# files/etc/uci-defaults/12-migrate-system.sh
#!/bin/sh 

. /lib/functions.sh 


config_load system || return 0 

