# files/etc/uci-defaults/20-autoconf.sh
#!/bin/sh 

uci set ffwizard.settings.firstboot=1
uci set ffwizard.settings.enabled=1
uci commit
