# files/etc/hotplug.d/ieee80211/99_dynack
#!/bin/ash

[ "$ACTION" = "add" ] || exit 0

phy=${DEVPATH##*/}
iw phy "$phy" set distance auto 

