# files/etc/hotplug.d/iface/00-cake-qdisc
#!/bin/sh

# FIXME https://bugs.lede-project.org/index.php?do=details&task_id=563
return

case "$ACTION" in
		ifup)
		grep -q '^sch_cake' '/proc/modules' && {
			command -V tc > /dev/null || return 
			tc qdisc replace dev "$DEVICE" root cake
		}
	;;
esac
# vim: set filetype=sh ai noet ts=4 sw=4 sts=4 :
