# files/etc/hotplug.d/iface/00-cake-qdisc
#!/bin/sh

case "$ACTION" in
        ifup)
		grep -q '^sch_cake' '/proc/modules' && {
			command -V tc > /dev/null || return 
			tc qdisc replace dev "$DEVICE" root cake
		}
	;;
esac
