# files/etc/hotplug.d/iface/99-registrator
#!/bin/sh 

[ "$INTERFACE" = "wan" ] && {

nn=$(uci -q get ffwizard.nodenumber)
[ -z "$nodenumber" ] && registrator autoconf 

}
